# fetcher/deploy/fetcher-api.service
[Unit]
Description=arxiv-firehose HTTP API (fetcher serve)
After=network-online.target
Wants=network-online.target

[Service]
# Installed as a *user* unit under ~/.config/systemd/user/ so no sudo is
# needed to deploy or restart. `loginctl enable-linger tower` (one-time)
# is what makes it survive logout / auto-start at boot.
Type=simple
WorkingDirectory=%h/apps/arxiv-firehose/fetcher
Environment=ARXIV_FIREHOSE_CACHE_DIR=/mnt/bertha/.cache/arxiv-firehose
Environment=ARXIV_FIREHOSE_LOG_DIR=/mnt/bertha/data/arxiv-firehose
# 0.0.0.0 binds every interface; the tailscale ACL + UFW are the perimeter.
ExecStart=%h/.local/bin/uv run fetcher serve \
    --data-dir /mnt/bertha/data/arxiv-firehose/data \
    --host 0.0.0.0 \
    --port 8087
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
